# crontab
0 * * * * /usr/local/bin/npm --prefix /app run start >> /var/log/cron.log 2>&1
