# crontab
0 */12 * * * /usr/local/bin/node /app/script.js >> /var/log/cron.log 2>&1
